# srt/supervisord.conf
[supervisord]
nodaemon=true
logfile=syslog
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel = DEBUG

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock
username=www
password=www

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:srt]
command = /opt/livesrt/bin/srt.sh
directory= /opt/livesrt
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:srtstats]
command = /opt/livesrt/bin/stats
directory= /opt/livesrt
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:noalbs]
command = /opt/livesrt/bin/noalbs
directory= /opt/noalbs
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0